# templates/etc/profile.d/composer.sh.j2
# {{ ansible_managed }}
export CGR_BIN_DIR="/usr/local/bin"
export COMPOSER_ALLOW_SUPERUSER="1"
export COMPOSER_BIN_DIR="/usr/local/bin"
export COMPOSER_HOME="/usr/local/share/composer"
