# templates/etc/profile.d/composer.sh.j2
# {{ ansible_managed }}
export COMPOSER_HOME=/usr/local/share/composer
export COMPOSER_BIN_DIR=/usr/local/bin
export CGR_BIN_DIR=/usr/local/share/composer
